# deployment/webp-editor.service
[Unit]
Description=WebP Editor FastAPI Application
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/path/to/webp-editor
Environment="PATH=/path/to/webp-editor/venv/bin"
ExecStart=/path/to/webp-editor/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 4

# Restart on failure
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/webp-editor/access.log
StandardError=append:/var/log/webp-editor/error.log

[Install]
WantedBy=multi-user.target 